# files/npm_lazy.conf
description "Daemon to cache / proxy npm module installs"

start on runlevel [2345]

respawn

script
  export HOME="/root"
  exec /usr/bin/npm_lazy -c /etc/npm_lazy-config.js
end script
